# infra/systemd/iaase-discover.service
[Unit]
Description=IAASE Marketplace - Discovery
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/dev/personal/IphoneScraperFlipper/scraper
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/bin/bash -lc "export NVM_DIR=\"$HOME/.nvm\"; [ -s \"$NVM_DIR/nvm.sh\" ] && . \"$NVM_DIR/nvm.sh\"; command -v nvm >/dev/null 2>&1 && nvm use --lts >/dev/null; %h/dev/personal/IphoneScraperFlipper/scraper/scripts/run_discover.sh"
ExecStartPost=/bin/bash -lc "export NVM_DIR=\"$HOME/.nvm\"; [ -s \"$NVM_DIR/nvm.sh\" ] && . \"$NVM_DIR/nvm.sh\"; command -v nvm >/dev/null 2>&1 && nvm use --lts >/dev/null; bash %h/dev/personal/IphoneScraperFlipper/scraper/scripts/run_compute_deals.sh"
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6

[Install]
WantedBy=default.target
